# etc/profile.d/voidbr-base-config-link-nano-to-ed.sh
#!/bin/sh

sudo ln -sf /usr/bin/nano /usr/local/bin/ed
